# framework/agent/scripts/flask.conf

Description "flask"
start on stopped rc RUNLEVEL=[2345]
respawn
exec python3 /home/vagrant/agent/agent.py

